# docker/supervisor/supervisord.ini
[supervisord]
childlogdir=/var/log/supervisord/
logfile=/var/log/supervisord/supervisord.log
pidfile=/run/supervisord/supervisord.pid
nodaemon=true

[program:nginx]
command=nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm81
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command=/usr/sbin/crond -f
autostart=true
autorestart=true
